# w1-watchdog.service
[Unit]
Description=1-Wire watchdog (power cycle on failure)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment="BUS_MASTER=/sys/bus/w1/devices/w1_bus_master1"
Environment="SENSOR_GLOB=28-*"
Environment="TEMP_FILE=temperature"
Environment="FAIL_THRESHOLD=3"
Environment="POWER_CYCLE_DELAY_SEC=5"
Environment="POWER_OFF_CMD=/usr/bin/wb-mqtt-publish -t /devices/wb-gpio/controls/5V_OUT -m 0"
Environment="POWER_ON_CMD=/usr/bin/wb-mqtt-publish -t /devices/wb-gpio/controls/5V_OUT -m 1"
Environment="LOG_TAG=w1-watchdog"
ExecStart=/opt/w1-monitor/w1-watchdog.sh
